# deploy/systemd/seneschal-scheduler-digest.timer
[Unit]
Description=Seneschal scheduled task example

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
